# ubifs-root/miwifi_rd03_firmware_ef0ee_1.0.47.bin/squashfs-root/etc/init.d/milog
#!/bin/sh /etc/rc.common
#
# This script is used for internal application traffic control
#

STOP=89

stop() {
	/usr/sbin/syslog_sync.sh

	return 0
}
